# examples/imraa/imraa.service
[Unit]
Description=imraa

[Service]
Type=oneshot
ExecStart=/home/longwei/Intel-IoT/mraa/build/examples/imraa/imraa